# build/supervisord.conf
[unix_http_server]
file=/run/supervisord.sock   ; (the path to the socket file)

[supervisord]
user=root
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0

[inet_http_server]         ; inet (TCP) server disabled by default
port=127.0.0.1:9001        ; ip_address:port specifier, *:port for all iface

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///run/supervisord.sock 


[program:erigon]
command=/usr/local/bin/erigon --datadir=/home/erigon --private.api.addr=localhost:9090 
#user=erigon
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:rpc]
command=/usr/local/bin/rpcdaemon --datadir=/home/erigon --txpool.api.addr=localhost:9090 --private.api.addr=localhost:9090 --http.api=eth,erigon,web3,net,debug,trace,txpool --ws
#user=erigon
#command=/app/start.sh
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
